# debian7/nvda-remote-server_1.5/DEBIAN/postrm
#!/bin/bash
# Post removal script.
 # This script runs after removing the package 

set -e
if test $1 = remove
then
if test -e /usr/share/NVDARemoteServer
then
rm -rf /usr/share/NVDARemoteServer
fi
update-rc.d NVDARemoteServer remove
fi
set +e
exit 0
